# supervisor.conf
[program:recognizer]
command=python /app/worker.py
process_name=%(process_num)s
numprocs=1
autorestart=true
startretries=50
redirect_stderr=true
stdout_logfile=/tmp/yp-worker.log
